# conf/shadowcaster-supervisor.conf
[program:shadowcaster]
command = /usr/bin/python /srv/run.py server
priority = 1